# src/bin/deploy/sing-box-orchestrator.service
[Unit]
Description=Sing-box client config and nodes orchestrator
After=network.target

[Service]
Type=simple
ExecStart=/etc/sing-box/bot
WorkingDirectory=/etc/sing-box
User=root
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target